# setup/blacklistmonitor.cron
0 9 * * 1     php /var/www/html/blacklistmonitor/service/cron.php -r weekly >>/var/log/blacklistmonitor.log 2>&1
@daily        php /var/www/html/blacklistmonitor/service/cron.php -r daily >>/var/log/blacklistmonitor.log 2>&1
0 */8 * * *   php /var/www/html/blacklistmonitor/service/cron.php -r 8hour >>/var/log/blacklistmonitor.log 2>&1
0 */2 * * *   php /var/www/html/blacklistmonitor/service/cron.php -r 2hour >>/var/log/blacklistmonitor.log 2>&1
@hourly       php /var/www/html/blacklistmonitor/service/cron.php -r 1hour >>/var/log/blacklistmonitor.log 2>&1

0 0 * * *     php /var/www/html/blacklistmonitor/service/cron.php -r blockListStats >>/var/log/blacklistmonitor.log 2>&1
@daily        php /var/www/html/blacklistmonitor/service/cron.php -r deleteOld >>/var/log/blacklistmonitor.log 2>&1
